# templates/default/dr-elephant.service.erb
[Unit]
Description = Drelephant server
After = syslog.target network.target remote-fs.target

[Service]
User = <%= node.drelephant.user %>
Environment=ELEPHANT_CONF_DIR= <%= node.drelephant.home %>/conf
PIDFile = /tmp/drelephant.pid
ExecStart = <%= node.drelephant.home %>/bin/start.sh
ExecStop = <%= node.drelephant.home %>/bin/stop.sh
LimitNOFILE = 20000

Type = forking

[Install]
WantedBy = multi-user.target
